# systemd/stuffemu.path
[Unit]
Description=Atari ST HW Emulator watch script
DefaultDependencies=no

[Install]
WantedBy=multi-user.target

[Path]
DirectoryNotEmpty=/dev/bus/usb/001/
